# cloud_native/week_2/task_1_packer_files/myapp.service
[Unit]
Description=Monolithic Application

[Service]
ExecStart=/home/packer/run_monolith.sh

[Install]
WantedBy=multi-user.target
